# deploy/all-in-one/supervisor/includes/disabled/embedded-redis.conf
[program:redis]
user=redis
command=/baserow/supervisor/wrapper.sh REDIS /usr/bin/redis-server /etc/redis/redis.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
autostart=true
priority=0
startretries=3
startsecs=30
